# llama-router.service
[Unit]
Description=llama-router - llama-server wrapper with TTL-based memory management
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=WORKDIR
ExecStart=WORKDIR/llama-router -preset WORKDIR/preset.ini -ttl TTL
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llama-router

[Install]
WantedBy=multi-user.target
